# conf.template/template.supervisor_admin.conf
[program:kdys-admin]
command=/home/job/kdys/kdys-api/env/bin/uwsgi /home/job/kdys/kdys-api/conf/uwsgi_admin.ini
user=job
directory=/home/job/kdys/kdys-api
autostart=true
autorestart=true
startsecs=3
startretries=10
stopasgroup=true
redirect_stderr=true
stdout_logfile=/data/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
